# dockerimages/metadatadns/supervisord.conf
[supervisord]
nodaemon=true

[program:bind9]
command=/usr/sbin/named -f -d 1 
stdout_events_enabled=true
stderr_events_enabled=true
